# daemons/process_alarm.timer
[Unit]
Description=Запуск монитора каждую минуту

[Timer]
OnUnitActiveSec=1min
Unit=process_alarm.service

[Install]
WantedBy=timers.target